# debian/deb/DEBIAN/postrm
#!/bin/bash

update-rc.d repetierlights remove
