# debian/deb/DEBIAN/postinst
#!/bin/bash

update-rc.d repetierlights defaults

/etc/init.d/repetierlights start

